# container/rootfs/etc/crontabs/root
# Run OSM updates every hour.
0 * * * * /app/bin/cron-hourly
