# container/rootfs/etc/crontabs/root
# Run OSM updates every hour.
0 * * * * cd /app; bin/treemap --osm-pull
0 * * * * cd /app; bin/treemap --upload-files
